# rpmbuild/SOURCES/dev50-0-1/etc/cron.d/rsnapshot
*/10 * * * * root /bin/nice /usr/bin/rsnapshot minutely
2 * * * * root /bin/nice /usr/bin/rsnapshot hourly
4 0 * * * root /bin/nice /usr/bin/rsnapshot daily
6 0 * * 0 root /bin/nice /usr/bin/rsnapshot weekly
